# systemd/fountain-dispatcher.service
[Unit]
Description=Codex FountainAI Deployment Dispatcher
After=network.target

[Service]
ExecStart=/usr/bin/python3 /srv/deploy/dispatcher_v2.py
WorkingDirectory=/srv/deploy
Restart=always
User=root

[Install]
WantedBy=multi-user.target
